# jenkins-k8s/templates/default/kubernetes-master.erb
#!/bin/bash
#
# kubernetes	This shell script takes care of starting and stopping kubernetes
#
# chkconfig: - 99 99
# description: kubernetes

### BEGIN INIT INFO
# Provides: kubernetes
# Required-Start: $network $local_fs $remote_fs
# Required-Stop: $network $local_fs $remote_fs
# Should-Start: $syslog $named flanneld docker
# Should-Stop: $syslog $named
# Short-Description: start and stop kubernetes
# Description: kubernetes
### END INIT INFO

# Source function library.
. /etc/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

prog=/usr/local/bin/hyperkube
lockfile=/var/lock/subsys/`basename $prog`
hostname=`hostname`
logfile=/var/log/kubernetes.log

CLUSTER_NAME=my-happy-cluster
ETCD_SERVERS="http://<%= @etcd_ba_account %>:<%= @etcd_ba_password %>@<%= @etcd_url %>:80"
CLUSTER_IP_RANGE="<%= @cluster_cidr %>"
MASTER="127.0.0.1:8080"


start() {
	#set IP range to etcd
	curl -L ${ETCD_SERVERS}/v2/keys/coreos.com/network/config \
		-XPUT -d value="{ \"Network\": \"${CLUSTER_IP_RANGE}\" }"


	# Start daemon.
	echo $"Starting apiserver: "
	daemon $prog apiserver \
	--service-cluster-ip-range=${CLUSTER_IP_RANGE} \
	--insecure-port=8080 \
	--secure-port=6443 \
	--basic-auth-file="<%= @ba_path %>" \
	--address=0.0.0.0 \
	--etcd_servers=${ETCD_SERVERS} \
	--cluster_name=${CLUSTER_NAME} \
	> ${logfile} 2>&1 &


	echo $"Starting controller-manager: "
	daemon $prog controller-manager \
	--master=${MASTER} \
	> ${logfile} 2>&1 &


	echo $"Starting scheduler: "
	daemon $prog scheduler \
	--master=${MASTER} \
	> ${logfile} 2>&1 &


	echo $"Starting proxy: "
	daemon $prog proxy \
	--master=${MASTER} \
	> ${logfile} 2>&1 &


	RETVAL=$?
	[ $RETVAL -eq 0 ] && touch $lockfile
	return $RETVAL
}

stop() {
	[ "$EUID" != "0" ] && exit 4
        echo -n $"Shutting down $prog: "
	killproc $prog
	RETVAL=$?
        echo
	[ $RETVAL -eq 0 ] && rm -f $lockfile
	return $RETVAL
}

# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  status)
	status $prog
	;;
  restart|force-reload)
	stop
	start
	;;
  try-restart|condrestart)
	if status $prog > /dev/null; then
	    stop
	    start
	fi
	;;
  reload)
	exit 3
	;;
  *)
	echo $"Usage: $0 {start|stop|status|restart|try-restart|force-reload}"
	exit 2
esac
